# modbus_server/crontab
*/30 * * * * python /usr/local/src/main.py
# Execute every 30 minutes